# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-849ce09c\x2d7d6d\x2d6b81\x2d2d98\x2db47769547880.target

[Mount]
What=/dev/disk/by-uuid/849ce09c-7d6d-6b81-2d98-b47769547880
Where=/
Type=ext4
Options=errors=remount-ro
